# tags/Gennaio2009_DifferentiTrigger/TDR_PC_SW/tkdev/GbatchInt/AMSdev/perl/parsejf.cron
#!/bin/bash
# parsejf daily 
/var/www/cgi-bin/mon/parsejf.o.cgi -c
exit 0


